# upstart.conf
start on runlevel [2345]

respawn
respawn limit 10 5 # max 10 times within 5 seconds

setuid ubuntu
chdir /home/ubuntu/pongout

exec node index.js
